# ansible/playbooks/templates/rc.local.j2
#!/bin/bash

ifconfig {{ secondNetIf_name }} {{ secondNetIf_ip }}/24 up

exit 0